# server/homepage.service
[Unit]
Description=Homepage static server (sven.guru)
After=network.target

[Service]
Type=simple
User=homepage
Group=homepage
ExecStart=/opt/homepage/static-server
WorkingDirectory=/opt/homepage
Restart=on-failure
RestartSec=5

Environment=PORT=8080
Environment=DOMAIN=sven.guru

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
