# server/homepage.service
[Unit]
Description=Homepage static server (sven.guru)
After=network.target

[Service]
Type=simple
User=homepage
Group=homepage
ExecStart=/opt/homepage/static-server
WorkingDirectory=/opt/homepage
Restart=on-failure
RestartSec=5

# Environment
Environment=ENABLE_HTTPS=true
Environment=DOMAIN=sven.guru
Environment=ACME_CONTACT_EMAIL=sven@sven.guru
Environment=ACME_STAGING=false
Environment=CERT_DIR=/opt/homepage

# Allow binding to ports 80 and 443
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/homepage
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
